# avrdisco.service
[Unit]
Description=AVRDisco Web Interface for AV Receiver Control
Documentation=https://github.com/lorton/AVRDisco-Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# IMPORTANT: Update these values for your setup
# If you used install_service.sh, these will be set automatically
# Otherwise, replace with your username and installation path:
User=pi
Group=pi
WorkingDirectory=/home/pi/AVRDisco-Web

# Update these paths to match your installation directory
Environment="PATH=/home/pi/AVRDisco-Web/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Load environment variables from .env file if it exists
EnvironmentFile=-/home/pi/AVRDisco-Web/.env

# Run the async version with hypercorn for production
# Update the path to match your installation directory
ExecStart=/home/pi/AVRDisco-Web/venv/bin/hypercorn async_app:app --bind 0.0.0.0:5000 --workers 2

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=avrdisco

[Install]
WantedBy=multi-user.target
